# supervisord.conf
[supervisord]
nodaemon=true

[program:kudu]
command=fastcgi-mono-server4 /applications=/:/opt/Kudu/ /socket=tcp:127.0.0.1:9000
priority=5

[program:nginx]
command=nginx -g "daemon off;"
priority=10